# docker-safe-sock.service
[Unit]
Description=Docker Safe Sock Proxy
After=docker.service network.target
Requires=docker.service

[Service]
Type=simple
EnvironmentFile=-/etc/default/docker-safe-sock
RuntimeDirectory=docker-safe
ExecStart=/usr/local/bin/docker-safe-sock
Restart=always
RestartSec=5
User=docker-safe
Group=docker-safe
# Security hardening
ProtectSystem=full
ProtectHome=true
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
